# docker/root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

# shellcheck disable=SC2086
# shellcheck shell=bash

umask 022
phpfpm_command="php-fpm7 -F"
echo "-> Starting php-fpm <-"
exec s6-setuidgid abc $phpfpm_command
